# deploy/systemd/iconia-admin.service
[Unit]
Description=ICONIA Admin (Next.js on :3000)
After=network-online.target iconia-server.service
Wants=network-online.target

[Service]
Type=simple
User=iconia
Group=iconia
# standalone 빌드는 자기 위치를 cwd 기준으로 .next/static / public 찾는다.
# build-and-upload.ps1 가 빌드 후 .next/standalone 으로 static/public 을 복사한다.
WorkingDirectory=/opt/iconia/admin/.next/standalone

EnvironmentFile=/etc/iconia.env
EnvironmentFile=-/etc/iconia.admin.env
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=NODE_OPTIONS=--max-old-space-size=1024
Environment=HOSTNAME=0.0.0.0

# next.config.mjs 의 output:'standalone' 산출물.
ExecStart=/usr/bin/node server.js
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
LimitNOFILE=65535

# systemd hardening.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/iconia/admin /var/log/iconia /tmp
PrivateTmp=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
PrivateDevices=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictSUIDSGID=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native

StandardOutput=journal
StandardError=journal
SyslogIdentifier=iconia-admin

[Install]
WantedBy=multi-user.target
